# data/systemd/sm.puri.Phosh.service.in
# This is a systemd user unit intended to be started with gnome-session.

[Unit]
Description=Phosh, a shell for mobile phones
Documentation=https://gitlab.gnome.org/World/Phosh/phosh

After=gnome-session-manager.target

Requisite=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session-initialized.target

StartLimitIntervalSec=15s
StartLimitBurst=3

OnFailure=gnome-session-shutdown.target
OnFailureJobMode=replace-irreversibly
CollectMode=inactive-or-failed
RefuseManualStart=on
RefuseManualStop=on

[Service]
Type=notify
ExecStart=@libexecdir@/phosh
Restart=on-failure
RestartSec=0ms
Slice=session.slice
